# supervisorconf/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
user=root
strip_ansi = true
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true           ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:auto_info_parser]
command=python3 /usr/src/app/auto_info_worker.py
numprocs=2
autostart=True
autorestart=True
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=logs/error_auto_info_worker.log
stderr_logfile=logs/log_auto_info_parser.log

[program:fines_parser]
command=python3 /usr/src/app/fines_worker.py
numprocs=2
autostart=True
autorestart=True
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=logs/error_fines_worker.log
stderr_logfile=logs/log_fines_worker.log
